# roles/appsuite/templates/etc/profile.d/open-xchange.sh.j2
{{ ansible_managed | comment }}

export PATH="/opt/open-xchange/sbin:$PATH"
